# cron.weekly/pfl
#!/bin/sh
exec nice /usr/bin/pfl >/dev/null
